# periodic/daily/999.updatedb
#!/bin/sh
#
# $FreeBSD$
#

/usr/libexec/locate.updatedb
